# ansible/roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd - highly-available key value store
Documentation=https://etcd.io/docs
After=network.target

[Service]
Type=notify
User=etcd
Group=etcd
ExecStart=/usr/local/bin/etcd --config-file /etc/etcd/etcd.conf.yml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
